# module_project1/workspace/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx 

[program:sshd]
command=/usr/sbin/sshd -D

